# package/DEBIAN/postinst
#!/bin/bash
set -e

echo "=========================================="
echo "   Настройка FridgeManager"
echo "=========================================="

chmod 755 /usr/bin/fridgemanager

# 🔧 ВАЖНО: Настраиваем локаль для правильной кодировки
echo "🔧 Настраиваем локаль и кодировку..."
export LANG="ru_RU.UTF-8"
export LC_ALL="ru_RU.UTF-8"

# Устанавливаем русскую локаль если нет
if ! locale -a | grep -q "ru_RU.utf8"; then
    echo "📝 Устанавливаем русскую локаль..."
    sudo locale-gen ru_RU.UTF-8
    sudo update-locale LANG=ru_RU.UTF-8
fi

if ! command -v psql &> /dev/null; then
    echo "❌ PostgreSQL не установлен!"
    exit 1
fi

echo "🔧 Запускаем PostgreSQL..."
systemctl start postgresql || true
systemctl enable postgresql || true

echo "⏳ Ожидаем запуск PostgreSQL..."
sleep 3

# 🔑 ВАЖНО: Определяем ТЕКУЩЕГО пользователя системы
CURRENT_USER=$(whoami)
echo "👤 Текущий пользователь системы: $CURRENT_USER"

# 🔑 СОЗДАЕМ ПОЛЬЗОВАТЕЛЯ PostgreSQL С ИМЕНЕМ ТЕКУЩЕГО СИСТЕМНОГО ПОЛЬЗОВАТЕЛЯ
echo "🔑 Создаем пользователя PostgreSQL: $CURRENT_USER"
sudo -u postgres createuser --superuser "$CURRENT_USER" 2>/dev/null && echo "✅ Пользователь $CURRENT_USER создан" || echo "⚠️ Пользователь $CURRENT_USER уже существует"

# 🔧 ПРОСТАЯ НАСТРОЙКА: Создаем базу и таблицы от имени postgres
echo "🗄️ Создаем базу данных..."
sudo -u postgres createdb fridgemanager 2>/dev/null && echo "✅ База данных создана" || echo "📁 База данных уже существует"

echo "📊 Создаем таблицу продуктов..."
sudo -u postgres psql -d fridgemanager -c "
CREATE TABLE IF NOT EXISTS products (
    id SERIAL PRIMARY KEY,
    name VARCHAR(100) UNIQUE NOT NULL,
    current_quantity INTEGER NOT NULL DEFAULT 0,
    norm_quantity INTEGER NOT NULL
);" 2>/dev/null && echo "✅ Таблица создана" || echo "📋 Таблица уже существует"

# Даем права ТЕКУЩЕМУ пользователю
sudo -u postgres psql -c "GRANT ALL PRIVILEGES ON DATABASE fridgemanager TO $CURRENT_USER;" 2>/dev/null && echo "✅ Права выданы пользователю $CURRENT_USER" || echo "⚠️ Права уже выданы"

# Даем права на таблицы
sudo -u postgres psql -d fridgemanager -c "GRANT ALL PRIVILEGES ON ALL TABLES IN SCHEMA public TO $CURRENT_USER;" 2>/dev/null && echo "✅ Права на таблицы выданы" || echo "⚠️ Права на таблицы уже выданы"
sudo -u postgres psql -d fridgemanager -c "GRANT ALL PRIVILEGES ON ALL SEQUENCES IN SCHEMA public TO $CURRENT_USER;" 2>/dev/null && echo "✅ Права на последовательности выданы" || echo "⚠️ Права на последовательности уже выданы"

echo "📝 Добавляем тестовые данные..."
sudo -u postgres psql -d fridgemanager -c "
INSERT INTO products (name, current_quantity, norm_quantity) VALUES
('Творог', 5, 10),
('Сыр', 12, 15), 
('Молоко', 18, 20),
('Яйца', 25, 30),
('Оливки', 3, 8)
ON CONFLICT (name) DO NOTHING;" 2>/dev/null && echo "✅ Данные добавлены" || echo "📊 Данные уже существуют"

# 🔧 ПРОСТАЯ ПРОВЕРКА: используем существующие настройки
echo "🔍 Проверяем подключения..."

# Проверяем от имени postgres (самый надежный способ)
if sudo -u postgres psql -d fridgemanager -c "SELECT COUNT(*) FROM products;" &>/dev/null; then
    COUNT=$(sudo -u postgres psql -d fridgemanager -t -c "SELECT COUNT(*) FROM products;" | tr -d ' ')
    echo "✅ Подключение postgres: РАБОТАЕТ"
    echo "📊 Продуктов в базе: $COUNT"
else
    echo "❌ Подключение postgres: НЕ РАБОТАЕТ"
fi

# Проверяем peer authentication для ТЕКУЩЕГО пользователя
if psql -d fridgemanager -c "SELECT current_user;" &>/dev/null; then
    USER=$(psql -d fridgemanager -t -c "SELECT current_user;" | tr -d ' ')
    echo "✅ Peer authentication: РАБОТАЕТ (пользователь: $USER)"
else
    echo "❌ Peer authentication: НЕ РАБОТАЕТ для $CURRENT_USER"
    
    # Проверяем с пользователем postgres
    if sudo -u postgres psql -d fridgemanager -c "SELECT current_user;" &>/dev/null; then
        echo "✅ Peer auth (postgres): РАБОТАЕТ"
    else
        echo "❌ Peer auth (postgres): НЕ РАБОТАЕТ"
    fi
fi

echo ""
echo "🎉 Установка завершена!"
echo "👤 Создан пользователь PostgreSQL: $CURRENT_USER"
echo "💡 Приложение попробует подключиться как: $CURRENT_USER"
echo "➡️ Запустите: fridgemanager"
echo "=========================================="