# package/DEBIAN/postinst
#!/bin/bash
set -e

echo "=========================================="
echo "   Настройка FridgeManager"
echo "=========================================="

chmod 755 /usr/bin/fridgemanager

if ! command -v psql &> /dev/null; then
    echo "❌ PostgreSQL не установлен!"
    exit 1
fi

echo "🔧 Запускаем PostgreSQL..."
systemctl start postgresql || true
systemctl enable postgresql || true

echo "⏳ Ожидаем запуск PostgreSQL..."
sleep 5

# 🔧 ВАЖНО: ПРИНУДИТЕЛЬНО настраиваем trust authentication
echo "🔧 НАСТРАИВАЕМ Trust Authentication..."
PG_HBA_FILE=$(find /etc/postgresql -name "pg_hba.conf" | head -1)

if [ -f "$PG_HBA_FILE" ]; then
    echo "📁 Найден файл: $PG_HBA_FILE"
    
    # Делаем резервную копию
    sudo cp "$PG_HBA_FILE" "${PG_HBA_FILE}.backup"
    echo "✅ Создана резервная копия"
    
    # 🔥 ПРИНУДИТЕЛЬНО заменяем настройки аутентификации
    echo "🔥 Заменяем настройки аутентификации на trust..."
    
    # Создаем временный файл с trust аутентификацией
    cat << 'EOF' | sudo tee "$PG_HBA_FILE" > /dev/null
# PostgreSQL Client Authentication Configuration File
# ===================================================
# FridgeManager - Trust Authentication Setup
local   all             postgres                                trust
local   all             all                                     trust
host    all             all             127.0.0.1/32            trust
host    all             all             ::1/128                 trust
local   replication     all                                     trust
host    replication     all             127.0.0.1/32            trust
host    replication     all             ::1/128                 trust
EOF

    echo "✅ Trust authentication настроена"
    
    # Перезапускаем PostgreSQL чтобы применить изменения
    sudo systemctl restart postgresql
    sleep 3
    echo "✅ PostgreSQL перезапущен"
else
    echo "❌ Файл pg_hba.conf не найден!"
    exit 1
fi

# 🔑 Создаем пользователя (теперь без пароля)
echo "🔑 Создаем пользователя БД..."
sudo -u postgres psql -c "CREATE USER fridgeuser;" 2>/dev/null && echo "✅ Пользователь создан: fridgeuser" || echo "⚠️ Пользователь уже существует"

# Создаем базу данных
echo "🗄️ Создаем базу данных..."
sudo -u postgres createdb fridgemanager 2>/dev/null && echo "✅ База данных создана" || echo "📁 База данных уже существует"

# Даем права
sudo -u postgres psql -c "GRANT ALL PRIVILEGES ON DATABASE fridgemanager TO fridgeuser;" 2>/dev/null && echo "✅ Права на БД выданы" || echo "⚠️ Права уже выданы"

echo "📊 Создаем таблицу продуктов..."
sudo -u postgres psql -d fridgemanager -c "
CREATE TABLE IF NOT EXISTS products (
    id SERIAL PRIMARY KEY,
    name VARCHAR(100) UNIQUE NOT NULL,
    current_quantity INTEGER NOT NULL DEFAULT 0,
    norm_quantity INTEGER NOT NULL
);" 2>/dev/null && echo "✅ Таблица создана" || echo "📋 Таблица уже существует"

# Даем права на таблицу
sudo -u postgres psql -d fridgemanager -c "GRANT ALL PRIVILEGES ON ALL TABLES IN SCHEMA public TO fridgeuser;" 2>/dev/null && echo "✅ Права на таблицы выданы" || echo "⚠️ Права уже выданы"
sudo -u postgres psql -d fridgemanager -c "GRANT ALL PRIVILEGES ON ALL SEQUENCES IN SCHEMA public TO fridgeuser;" 2>/dev/null && echo "✅ Права на последовательности выданы" || echo "⚠️ Права уже выданы"

echo "📝 Добавляем тестовые данные..."
sudo -u postgres psql -d fridgemanager -c "
INSERT INTO products (name, current_quantity, norm_quantity) VALUES
('Творог', 5, 10),
('Сыр', 12, 15), 
('Молоко', 18, 20),
('Яйца', 25, 30),
('Оливки', 3, 8)
ON CONFLICT (name) DO NOTHING;" 2>/dev/null && echo "✅ Данные добавлены" || echo "📊 Данные уже существуют"

echo "🔍 Проверяем подключения..."

# Проверяем подключение БЕЗ пароля (trust auth)
if psql -h localhost -p 5432 -U fridgeuser -d fridgemanager -c "SELECT current_user, COUNT(*) FROM products;" &>/dev/null; then
    echo "✅ Trust authentication: РАБОТАЕТ"
    PRODUCT_COUNT=$(psql -h localhost -p 5432 -U fridgeuser -d fridgemanager -t -c "SELECT COUNT(*) FROM products;" | tr -d ' ')
    echo "📊 Продуктов в базе: $PRODUCT_COUNT"
else
    echo "❌ Trust authentication: НЕ РАБОТАЕТ"
    echo "Проверяем настройки pg_hba.conf:"
    sudo cat "$PG_HBA_FILE"
    exit 1
fi

echo ""
echo "🎉 Установка завершена!"
echo "🔐 Trust authentication активирована"
echo "   Пользователь: fridgeuser"
echo "   Пароль: НЕ ТРЕБУЕТСЯ"
echo "   База данных: fridgemanager"
echo "   Хост: localhost"
echo "   Порт: 5432"
echo "➡️ Запустите: fridgemanager"
echo "=========================================="