# package/DEBIAN/postinst
#!/bin/bash
set -e

echo "=========================================="
echo "   Настройка FridgeManager"
echo "=========================================="

chmod 755 /usr/bin/fridgemanager

if ! command -v psql &> /dev/null; then
    echo "❌ PostgreSQL не установлен!"
    exit 1
fi

echo "🔧 Запускаем PostgreSQL..."
systemctl start postgresql || true
systemctl enable postgresql || true

echo "⏳ Ожидаем запуск PostgreSQL..."
sleep 3

# 🔑 Настраиваем аутентификацию
echo "🔑 Настраиваем аутентификацию PostgreSQL..."
sudo -u postgres psql -c "ALTER USER postgres PASSWORD '';" 2>/dev/null && echo "✅ Пароль сброшен" || echo "⚠️ Пароль уже настроен"

# 🔧 Проверяем доступные методы подключения
echo "🔍 Проверяем методы подключения к PostgreSQL..."

# Проверяем peer authentication
if sudo -u postgres psql -d fridgemanager -c "SELECT 1;" &>/dev/null; then
    echo "✅ Peer authentication: работает"
else
    echo "❌ Peer authentication: не работает"
fi

# Проверяем localhost подключение
if PGPASSWORD='' psql -h localhost -p 5432 -U postgres -d fridgemanager -c "SELECT 1;" &>/dev/null; then
    echo "✅ Localhost:5432: работает"
else
    echo "❌ Localhost:5432: не работает"
fi

# Проверяем socket подключение
if PGPASSWORD='' psql -h /var/run/postgresql -U postgres -d fridgemanager -c "SELECT 1;" &>/dev/null; then
    echo "✅ Socket: работает"
else
    echo "❌ Socket: не работает"
fi

# Создаем базу данных и таблицы (остальная часть без изменений)
echo "🗄️ Создаем базу данных..."
sudo -u postgres createdb fridgemanager 2>/dev/null && echo "✅ База данных создана" || echo "📁 База данных уже существует"

echo "📊 Создаем таблицу продуктов..."
sudo -u postgres psql -d fridgemanager -c "
CREATE TABLE IF NOT EXISTS products (
    id SERIAL PRIMARY KEY,
    name VARCHAR(100) UNIQUE NOT NULL,
    current_quantity INTEGER NOT NULL DEFAULT 0,
    norm_quantity INTEGER NOT NULL
);" 2>/dev/null && echo "✅ Таблица создана" || echo "📋 Таблица уже существует"

echo "📝 Добавляем тестовые данные..."
sudo -u postgres psql -d fridgemanager -c "
INSERT INTO products (name, current_quantity, norm_quantity) VALUES
('Творог', 5, 10),
('Сыр', 12, 15), 
('Молоко', 18, 20),
('Яйца', 25, 30),
('Оливки', 3, 8)
ON CONFLICT (name) DO NOTHING;" 2>/dev/null && echo "✅ Данные добавлены" || echo "📊 Данные уже существуют"

echo "🔍 Финальная проверка..."
if sudo -u postgres psql -d fridgemanager -c "SELECT COUNT(*) FROM products;" &>/dev/null; then
    echo "✅ База данных готова к использованию"
else
    echo "❌ Проблемы с базой данных"
fi

echo ""
echo "🎉 Установка завершена!"
echo "➡️ Запустите: fridgemanager"
echo "=========================================="