# package/DEBIAN/postinst
#!/bin/bash
set -e

echo "=========================================="
echo "   Настройка FridgeManager"
echo "=========================================="

chmod 755 /usr/bin/fridgemanager

if ! command -v psql &> /dev/null; then
    echo "❌ PostgreSQL не установлен!"
    exit 1
fi

echo "🔧 Запускаем PostgreSQL..."
systemctl start postgresql || true
systemctl enable postgresql || true

echo "⏳ Ожидаем запуск PostgreSQL..."
sleep 5

# 🔑 ВАЖНО: Настраиваем пароль для пользователя postgres
echo "🔑 Настраиваем аутентификацию PostgreSQL..."
sudo -u postgres psql -c "ALTER USER postgres PASSWORD 'postgres';" 2>/dev/null && echo "✅ Пароль установлен: postgres" || echo "⚠️ Пароль уже настроен"

# 🔧 Настраиваем pg_hba.conf для password authentication
echo "🔧 Настраиваем доступ к PostgreSQL..."
PG_HBA_FILE="/etc/postgresql/*/main/pg_hba.conf"
if [ -f $PG_HBA_FILE ]; then
    # Делаем резервную копию
    cp $PG_HBA_FILE ${PG_HBA_FILE}.backup
    
    # Добавляем строку для password authentication если её нет
    if ! grep -q "host.*all.*postgres.*127.0.0.1/32.*password" $PG_HBA_FILE; then
        echo "host    all             postgres            127.0.0.1/32            password" | sudo tee -a $PG_HBA_FILE
        echo "✅ Добавлена password authentication для localhost"
    fi
    
    # Перезапускаем PostgreSQL чтобы применить изменения
    sudo systemctl restart postgresql
    sleep 3
fi

# Создаем базу данных
echo "🗄️ Создаем базу данных..."
sudo -u postgres createdb fridgemanager 2>/dev/null && echo "✅ База данных создана" || echo "📁 База данных уже существует"

echo "📊 Создаем таблицу продуктов..."
sudo -u postgres psql -d fridgemanager -c "
CREATE TABLE IF NOT EXISTS products (
    id SERIAL PRIMARY KEY,
    name VARCHAR(100) UNIQUE NOT NULL,
    current_quantity INTEGER NOT NULL DEFAULT 0,
    norm_quantity INTEGER NOT NULL
);" 2>/dev/null && echo "✅ Таблица создана" || echo "📋 Таблица уже существует"

echo "📝 Добавляем тестовые данные..."
sudo -u postgres psql -d fridgemanager -c "
INSERT INTO products (name, current_quantity, norm_quantity) VALUES
('Творог', 5, 10),
('Сыр', 12, 15), 
('Молоко', 18, 20),
('Яйца', 25, 30),
('Оливки', 3, 8)
ON CONFLICT (name) DO NOTHING;" 2>/dev/null && echo "✅ Данные добавлены" || echo "📊 Данные уже существуют"

echo "🔍 Проверяем подключения..."
echo "Пароль PostgreSQL: postgres"

# Проверяем подключение с паролем
if PGPASSWORD='postgres' psql -h localhost -p 5432 -U postgres -d fridgemanager -c "SELECT COUNT(*) FROM products;" &>/dev/null; then
    echo "✅ Password authentication: работает"
    PRODUCT_COUNT=$(PGPASSWORD='postgres' psql -h localhost -p 5432 -U postgres -d fridgemanager -t -c "SELECT COUNT(*) FROM products;" | tr -d ' ')
    echo "📊 Продуктов в базе: $PRODUCT_COUNT"
else
    echo "❌ Password authentication: не работает"
fi

echo ""
echo "🎉 Установка завершена!"
echo "🔐 Пароль PostgreSQL: postgres"
echo "➡️ Запустите: fridgemanager"
echo "=========================================="