# package/DEBIAN/postinst
#!/bin/bash
set -e

echo "=========================================="
echo "   Настройка FridgeManager"
echo "=========================================="

chmod 755 /usr/bin/fridgemanager

if ! command -v psql &> /dev/null; then
    echo "❌ PostgreSQL не установлен!"
    exit 1
fi

echo "🔧 Запускаем PostgreSQL..."
systemctl start postgresql || true
systemctl enable postgresql || true

echo "⏳ Ожидаем запуск PostgreSQL..."
sleep 5

# 🔑 ВАЖНО: Создаем нового пользователя с известным паролем
echo "🔑 Создаем пользователя БД..."
sudo -u postgres psql -c "CREATE USER fridgeuser WITH PASSWORD 'fridge123';" 2>/dev/null && echo "✅ Пользователь создан" || echo "⚠️ Пользователь уже существует"

# 🔑 Устанавливаем пароль если пользователь уже есть
sudo -u postgres psql -c "ALTER USER fridgeuser PASSWORD 'fridge123';" 2>/dev/null && echo "✅ Пароль установлен" || echo "⚠️ Не удалось установить пароль"

# 🔑 Даем права на создание БД
sudo -u postgres psql -c "ALTER USER fridgeuser CREATEDB;" 2>/dev/null && echo "✅ Права выданы" || echo "⚠️ Права уже выданы"

# Создаем базу данных от имени нового пользователя
echo "🗄️ Создаем базу данных..."
sudo -u postgres createdb fridgemanager 2>/dev/null && echo "✅ База данных создана" || echo "📁 База данных уже существует"

# 🔑 Даем права пользователю на базу данных
sudo -u postgres psql -c "GRANT ALL PRIVILEGES ON DATABASE fridgemanager TO fridgeuser;" 2>/dev/null && echo "✅ Права на БД выданы" || echo "⚠️ Права уже выданы"

echo "📊 Создаем таблицу продуктов..."
sudo -u postgres psql -d fridgemanager -c "
CREATE TABLE IF NOT EXISTS products (
    id SERIAL PRIMARY KEY,
    name VARCHAR(100) UNIQUE NOT NULL,
    current_quantity INTEGER NOT NULL DEFAULT 0,
    norm_quantity INTEGER NOT NULL
);" 2>/dev/null && echo "✅ Таблица создана" || echo "📋 Таблица уже существует"

# 🔑 Даем права на таблицу
sudo -u postgres psql -d fridgemanager -c "GRANT ALL PRIVILEGES ON TABLE products TO fridgeuser;" 2>/dev/null && echo "✅ Права на таблицу выданы" || echo "⚠️ Права уже выданы"
sudo -u postgres psql -d fridgemanager -c "GRANT ALL PRIVILEGES ON SEQUENCE products_id_seq TO fridgeuser;" 2>/dev/null && echo "✅ Права на последовательность выданы" || echo "⚠️ Права уже выданы"

echo "📝 Добавляем тестовые данные..."
sudo -u postgres psql -d fridgemanager -c "
INSERT INTO products (name, current_quantity, norm_quantity) VALUES
('Творог', 5, 10),
('Сыр', 12, 15), 
('Молоко', 18, 20),
('Яйца', 25, 30),
('Оливки', 3, 8)
ON CONFLICT (name) DO NOTHING;" 2>/dev/null && echo "✅ Данные добавлены" || echo "📊 Данные уже существуют"

echo "🔍 Проверяем подключение с новым пользователем..."
if PGPASSWORD='fridge123' psql -h localhost -p 5432 -U fridgeuser -d fridgemanager -c "SELECT COUNT(*) FROM products;" &>/dev/null; then
    echo "✅ Подключение с fridgeuser: работает"
    PRODUCT_COUNT=$(PGPASSWORD='fridge123' psql -h localhost -p 5432 -U fridgeuser -d fridgemanager -t -c "SELECT COUNT(*) FROM products;" | tr -d ' ')
    echo "📊 Продуктов в базе: $PRODUCT_COUNT"
else
    echo "❌ Подключение с fridgeuser: не работает"
    echo "Пробуем альтернативные методы..."
    
    # Проверяем peer authentication для postgres
    if sudo -u postgres psql -d fridgemanager -c "SELECT COUNT(*) FROM products;" &>/dev/null; then
        echo "✅ Peer auth (postgres): работает"
    else
        echo "❌ Peer auth (postgres): не работает"
    fi
fi

echo ""
echo "🎉 Установка завершена!"
echo "🔐 Данные для подключения:"
echo "   Пользователь: fridgeuser"
echo "   Пароль: fridge123"
echo "   База данных: fridgemanager"
echo "   Хост: localhost"
echo "   Порт: 5432"
echo "➡️ Запустите: fridgemanager"
echo "=========================================="