# package/DEBIAN/postinst
#!/bin/bash
set -e

echo "=========================================="
echo "   Настройка FridgeManager"
echo "=========================================="

chmod 755 /usr/bin/fridgemanager

if ! command -v psql &> /dev/null; then
    echo "❌ PostgreSQL не установлен!"
    exit 1
fi

echo "🔧 Запускаем PostgreSQL..."
systemctl start postgresql || true
systemctl enable postgresql || true

echo "⏳ Ожидаем запуск PostgreSQL..."
sleep 5

# 🔧 ВАЖНО: Настраиваем pg_hba.conf для password authentication
echo "🔧 Настраиваем доступ к PostgreSQL..."
PG_HBA_FILE=$(find /etc/postgresql -name "pg_hba.conf" | head -1)

if [ -f "$PG_HBA_FILE" ]; then
    echo "📁 Найден файл: $PG_HBA_FILE"
    
    # Делаем резервную копию
    sudo cp "$PG_HBA_FILE" "${PG_HBA_FILE}.backup"
    echo "✅ Создана резервная копия"
    
    # Добавляем строки для password authentication если их нет
    if ! grep -q "host.*all.*all.*127.0.0.1/32.*password" "$PG_HBA_FILE"; then
        echo "# FridgeManager - password authentication" | sudo tee -a "$PG_HBA_FILE"
        echo "host    all             all             127.0.0.1/32            password" | sudo tee -a "$PG_HBA_FILE"
        echo "host    all             all             ::1/128                 password" | sudo tee -a "$PG_HBA_FILE"
        echo "✅ Добавлена password authentication для localhost"
    else
        echo "✅ Password authentication уже настроена"
    fi
    
    # Перезапускаем PostgreSQL чтобы применить изменения
    sudo systemctl restart postgresql
    sleep 3
    echo "✅ PostgreSQL перезапущен"
else
    echo "❌ Файл pg_hba.conf не найден!"
fi

# 🔑 Создаем нового пользователя
echo "🔑 Создаем пользователя БД..."
sudo -u postgres psql -c "CREATE USER fridgeuser WITH PASSWORD 'fridge123';" 2>/dev/null && echo "✅ Пользователь создан: fridgeuser" || echo "⚠️ Пользователь уже существует"

# Устанавливаем пароль
sudo -u postgres psql -c "ALTER USER fridgeuser PASSWORD 'fridge123';" 2>/dev/null && echo "✅ Пароль установлен: fridge123" || echo "⚠️ Не удалось установить пароль"

# Создаем базу данных
echo "🗄️ Создаем базу данных..."
sudo -u postgres createdb fridgemanager 2>/dev/null && echo "✅ База данных создана" || echo "📁 База данных уже существует"

# Даем права
sudo -u postgres psql -c "GRANT ALL PRIVILEGES ON DATABASE fridgemanager TO fridgeuser;" 2>/dev/null && echo "✅ Права на БД выданы" || echo "⚠️ Права уже выданы"

echo "📊 Создаем таблицу продуктов..."
sudo -u postgres psql -d fridgemanager -c "
CREATE TABLE IF NOT EXISTS products (
    id SERIAL PRIMARY KEY,
    name VARCHAR(100) UNIQUE NOT NULL,
    current_quantity INTEGER NOT NULL DEFAULT 0,
    norm_quantity INTEGER NOT NULL
);" 2>/dev/null && echo "✅ Таблица создана" || echo "📋 Таблица уже существует"

# Даем права на таблицу
sudo -u postgres psql -d fridgemanager -c "GRANT ALL PRIVILEGES ON ALL TABLES IN SCHEMA public TO fridgeuser;" 2>/dev/null && echo "✅ Права на таблицы выданы" || echo "⚠️ Права уже выданы"
sudo -u postgres psql -d fridgemanager -c "GRANT ALL PRIVILEGES ON ALL SEQUENCES IN SCHEMA public TO fridgeuser;" 2>/dev/null && echo "✅ Права на последовательности выданы" || echo "⚠️ Права уже выданы"

echo "📝 Добавляем тестовые данные..."
sudo -u postgres psql -d fridgemanager -c "
INSERT INTO products (name, current_quantity, norm_quantity) VALUES
('Творог', 5, 10),
('Сыр', 12, 15), 
('Молоко', 18, 20),
('Яйца', 25, 30),
('Оливки', 3, 8)
ON CONFLICT (name) DO NOTHING;" 2>/dev/null && echo "✅ Данные добавлены" || echo "📊 Данные уже существуют"

echo "🔍 Проверяем подключения..."

# Проверяем подключение с новым пользователем
if PGPASSWORD='fridge123' psql -h localhost -p 5432 -U fridgeuser -d fridgemanager -c "SELECT current_user;" &>/dev/null; then
    echo "✅ Подключение с fridgeuser: РАБОТАЕТ"
    PRODUCT_COUNT=$(PGPASSWORD='fridge123' psql -h localhost -p 5432 -U fridgeuser -d fridgemanager -t -c "SELECT COUNT(*) FROM products;" | tr -d ' ')
    echo "📊 Продуктов в базе: $PRODUCT_COUNT"
else
    echo "❌ Подключение с fridgeuser: НЕ РАБОТАЕТ"
    
    # Показываем текущие настройки pg_hba
    echo "📋 Текущие настройки аутентификации:"
    sudo grep -E "(host|local)" "$PG_HBA_FILE" | head -10
fi

echo ""
echo "🎉 Установка завершена!"
echo "🔐 Данные для подключения:"
echo "   Пользователь: fridgeuser"
echo "   Пароль: fridge123" 
echo "   База данных: fridgemanager"
echo "   Хост: localhost"
echo "   Порт: 5432"
echo "➡️ Запустите: fridgemanager"
echo "=========================================="