# package/DEBIAN/postinst
#!/bin/bash
set -e

echo "=========================================="
echo "   Настройка FridgeManager"
echo "=========================================="

chmod 755 /usr/bin/fridgemanager

if ! command -v psql &> /dev/null; then
    echo "❌ PostgreSQL не установлен!"
    exit 1
fi

echo "🔧 Запускаем PostgreSQL..."
systemctl start postgresql || true
systemctl enable postgresql || true

echo "⏳ Ожидаем запуск PostgreSQL..."
sleep 5

# 🔑 СОЗДАЕМ ПОЛЬЗОВАТЕЛЯ ДЛЯ PEER AUTHENTICATION
echo "🔑 Настраиваем пользователей PostgreSQL..."
CURRENT_USER=$(whoami)

# Создаем пользователя с таким же именем как текущий системный пользователь
sudo -u postgres createuser --superuser "$CURRENT_USER" 2>/dev/null && echo "✅ Создан пользователь: $CURRENT_USER" || echo "⚠️ Пользователь $CURRENT_USER уже существует"

# Сбрасываем пароль для postgres
sudo -u postgres psql -c "ALTER USER postgres PASSWORD '';" 2>/dev/null && echo "✅ Пароль postgres сброшен" || echo "⚠️ Пароль postgres уже настроен"

# Создаем базу данных
echo "🗄️ Создаем базу данных..."
sudo -u postgres createdb fridgemanager 2>/dev/null && echo "✅ База данных создана" || echo "📁 База данных уже существует"

echo "📊 Создаем таблицу продуктов..."
sudo -u postgres psql -d fridgemanager -c "
CREATE TABLE IF NOT EXISTS products (
    id SERIAL PRIMARY KEY,
    name VARCHAR(100) UNIQUE NOT NULL,
    current_quantity INTEGER NOT NULL DEFAULT 0,
    norm_quantity INTEGER NOT NULL
);" 2>/dev/null && echo "✅ Таблица создана" || echo "📋 Таблица уже существует"

echo "📝 Добавляем тестовые данные..."
sudo -u postgres psql -d fridgemanager -c "
INSERT INTO products (name, current_quantity, norm_quantity) VALUES
('Творог', 5, 10),
('Сыр', 12, 15), 
('Молоко', 18, 20),
('Яйца', 25, 30),
('Оливки', 3, 8)
ON CONFLICT (name) DO NOTHING;" 2>/dev/null && echo "✅ Данные добавлены" || echo "📊 Данные уже существуют"

echo "🔍 Проверяем подключения..."
echo "Текущий пользователь системы: $CURRENT_USER"

# Проверяем разные методы подключения
if psql -d fridgemanager -c "SELECT current_user;" &>/dev/null; then
    echo "✅ Peer authentication: работает"
else
    echo "❌ Peer authentication: не работает"
fi

if sudo -u postgres psql -d fridgemanager -c "SELECT current_user;" &>/dev/null; then
    echo "✅ Peer auth (postgres): работает"
else
    echo "❌ Peer auth (postgres): не работает"
fi

echo ""
echo "🎉 Установка завершена!"
echo "➡️ Запустите: fridgemanager"
echo "=========================================="